# src/modules/whitenode/filesystem/root_init/etc/systemd/system/whitenode_first_boot.service
[Unit]
Description=on first boot generate RPC user and pass for Whitecoind
[Service]
ExecStart=/home/pi/scripts/firstboot
ExecStop=
[Install]
WantedBy=multi-user.target